# deploy/garagepro.service.d/override.conf
# Drop-in for garagepro.service — install to:
#   /etc/systemd/system/garagepro.service.d/override.conf
# then: systemctl daemon-reload && systemctl restart garagepro
#
# Allows the moderator one-click updater to escalate via the scoped sudo rule
# (/etc/sudoers.d/garagepro-update -> /usr/local/bin/garagepro-update only).
# Without this, NoNewPrivileges=yes in the hardened unit blocks sudo inside the
# service sandbox.
[Service]
NoNewPrivileges=no
